# miniconda/etc/profile.d/conda.csh
setenv CONDA_EXE "`cygpath 'F:\My Personals\MedicalChatbot\miniconda\Scripts\conda.exe'`";
setenv _CONDA_EXE "`cygpath 'F:\My Personals\MedicalChatbot\miniconda\Scripts\conda.exe'`";
unsetenv _CE_M;
unsetenv _CE_CONDA;
setenv CONDA_PYTHON_EXE "`cygpath 'F:\My Personals\MedicalChatbot\miniconda\python.exe'`";
setenv _CONDA_ROOT "`cygpath 'F:\My Personals\MedicalChatbot\miniconda'`";

source "`cygpath 'F:\My Personals\MedicalChatbot\miniconda\Lib\site-packages\conda\shell\etc\profile.d\conda.csh'`"